# pokemon bot claude version/pokemon-bot.service
[Unit]
Description=Pokemon Center Drop Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pokemon-bot
ExecStart=/usr/bin/python3 /home/pi/pokemon-bot/bot.py
Restart=always
RestartSec=10
StandardOutput=append:/home/pi/pokemon-bot/bot.log
StandardError=append:/home/pi/pokemon-bot/bot.log

[Install]
WantedBy=multi-user.target
